# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nestjs]
command=node /app/backend/dist/main.js
autostart=true
autorestart=true
stdout_logfile=/app/logs/nestjs.log
stderr_logfile=/app/logs/nestjs-error.log
priority=10

[program:python]
command=python3.11 -m uvicorn ai-service.main:app --host 0.0.0.0 --port 8000 --workers 2
autostart=true
autorestart=true
stdout_logfile=/app/logs/python.log
stderr_logfile=/app/logs/python-error.log
priority=20
